# .config/autostart/qemu.desktop
[Desktop Entry]
Name=QEMU Session
Exec=/usr/bin/virsh -c qemu:///session
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
